# debmeta/DEBIAN/prerm
#!/bin/sh -e

systemctl stop sample-todo.service
systemctl disable sample-todo.service
